# com.walmart.fashion-portal-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- E2E Fashion Transformation Portal — Daily OPIF Sync -->
  <!-- Runs daily-sync.py every day at 8:00 AM -->
  <key>Label</key>
  <string>com.walmart.fashion-portal-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/e0c0lzr/Downloads/fashion-portal/daily-sync.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal</string>

  <key>StandardOutPath</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal/sync.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal/sync.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/e0c0lzr</string>
  </dict>

  <!-- Retry if machine was asleep at scheduled time -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Keep alive if it crashes -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>