# com.walmart.fashion-portal-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- E2E Fashion Portal — Local OPIF Proxy (runs persistently) -->
  <!-- Allows the portal dashboard to fetch live OPIF data via localhost:7465 -->
  <key>Label</key>
  <string>com.walmart.fashion-portal-proxy</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/e0c0lzr/Downloads/fashion-portal/.venv/bin/python3</string>
    <string>/Users/e0c0lzr/Downloads/fashion-portal/opif-proxy.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal</string>

  <key>StandardOutPath</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal/proxy.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/e0c0lzr/Downloads/fashion-portal/proxy.log</string>

  <!-- Restart automatically if it crashes -->
  <key>KeepAlive</key>
  <true/>

  <!-- Start when the plist is loaded -->
  <key>RunAtLoad</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/e0c0lzr</string>
  </dict>
</dict>
</plist>